# scripts/databridge-celery.service
[Unit]
Description=DataBridge Celery Worker
After=network.target redis.service

[Service]
Type=simple
User=nilesh.kute
WorkingDirectory=/opt/webapp/databridge-pipeline
ExecStart=/opt/webapp/databridge-pipeline/backend/.venv/bin/celery -A backend.app.core.celery_app worker -l info -Q scanning,transfer,notifications --concurrency=4
Restart=always
RestartSec=5
Environment=PATH=/opt/webapp/databridge-pipeline/backend/.venv/bin

[Install]
WantedBy=multi-user.target
