# scripts/databridge.service
[Unit]
Description=DataBridge Pipeline Server
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=databridge
WorkingDirectory=/opt/databridge-pipeline
ExecStart=/opt/databridge-pipeline/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=5
Environment=PATH=/opt/databridge-pipeline/backend/.venv/bin

[Install]
WantedBy=multi-user.target
